# Docker/composer/crontab
0 */12 * * * php /usr/share/nginx/pinko/bin/console app:instagram:sync
